# dashboard/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock
chmod=0775

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[inet_http_server]
port=0.0.0.0:9001
username=admin
password=admin123

[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/tmp/supervisord.pid
childlogdir=/var/log/supervisor

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:laravel-queue]
command=php /var/www/html/artisan queue:work
autostart=true
autorestart=true

[program:laravel-queue_long]
command=php /var/www/html/artisan queue:work --queue=long
autostart=true
autorestart=true

[program:laravel-queue-emails]
command=php /var/www/html/artisan queue:work redis --queue=emails --sleep=3 --tries=3 --timeout=90
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/queue-emails.log

[program:laravel-schedule]
command=php /var/www/html/artisan schedule:work
autostart=true
autorestart=true
